# cosmic-hdr.service
[Unit]
Description=HDR10 color pipeline (BT.2020 + PQ + metadata)
Documentation=https://github.com/Sigmachan/cosmic-hdr
# Start after the display manager brings up a graphical session.
# cosmic-hdr retries drmSetMaster internally (up to 6 s) so no sleep needed.
After=display-manager.service
Wants=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cosmic-hdr
ExecStop=/usr/local/bin/cosmic-hdr reset
RemainAfterExit=yes
# Retry if compositor is still initialising when the service fires
Restart=on-failure
RestartSec=3s
StartLimitBurst=4
StartLimitIntervalSec=30s

[Install]
WantedBy=multi-user.target
